# src/bios-agent-nut-configurator.service.in
[Unit]
Description=bios-agent-nut-configurator service
Requires=malamute.service
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/bios-agent-nut-configurator
Restart=always

[Install]
WantedBy=bios.target
