# src/bios-agent-nut-configurator.service.in
[Unit]
Description=bios-agent-nut-configurator service
After=malamute.target

[Service]
Type=simple
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/bios-agent-nut-configurator

[Install]
WantedBy=multi-user.target
