# src/fty-nut-configuration.service.in
# This is a skeleton created by zproject.
# You can add hand-written code here.

[Unit]
Description=fty-nut-configuration service
Requires=network.target malamute.service
After=network.target malamute.service
PartOf=bios.target

# This implies starting-after/stopping-before fty-db.target as well
After=fty-asset.service
Requires=fty-asset.service
BindsTo=fty-asset.service

After=fty-security-wallet.service
Requires=fty-security-wallet.service
BindsTo=fty-security-wallet.service

[Service]
Type=simple
User=discovery-monitoring-daemon
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
Environment='SYSTEMD_UNIT_FULLNAME=%n'
EnvironmentFile=/etc/default/bios-db-rw
ExecStart=@prefix@/bin/fty-nut-configuration -c @sysconfdir@/@PACKAGE@/fty-nut-configuration.cfg

[Install]
WantedBy=bios.target
