# src/fty-nut.service.in
# This is a skeleton created by zproject.
# You can add hand-written code here.

[Unit]
Description=fty-nut service
Requires=network.target malamute.service
After=network.target malamute.service
Before=fty-asset.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__agent-nut.service.conf
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__agent-nut.service.conf
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-nut --mapping-file @datadir@/fty-nut/mapping.conf --state-file '/var/lib/fty/fty-nut/state_file'

[Install]
WantedBy=bios.target
