# src/fty-nut.service.in
[Unit]
Description=fty-nut service
Requires=network.target malamute.service
After=network.target malamute.service
Before=fty-asset.service
PartOf=bios.target


[Service]
Type=simple
User=bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-nut --mapping-file @datadir@/@PACKAGE@/mapping.conf --state-file '/var/lib/fty/nut/state_file'
Restart=always

[Install]
WantedBy=bios.target
