# src/fty-nut.service.in
[Unit]
Description=fty-nut service
Requires=network.target malamute.service
After=network.target malamute.service
Before=fty-asset.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__agent-nut.service.conf
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__agent-nut.service.conf
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
Environment="prefix=@prefix@"
# TODO: Change later to parametrised names and un-legacy bios=>fty (and remove explicit service files above):
#ExecStart=@prefix@/bin/fty-nut --mapping-file @datadir@/@PACKAGE@/mapping.conf --state-file '/var/lib/fty/nut/state_file'
ExecStart=@prefix@/bin/fty-nut --mapping-file @datadir@/agent-nut/mapping.conf --state-file '/var/lib/bios/nut/state_file'
Restart=always

[Install]
WantedBy=bios.target
